# infrastructure/monitoring/monitoring-server/systemd/grafana-server.service.d/override.conf
[Service]
Environment=GF_SERVER_HTTP_ADDR=0.0.0.0
